# readonly-root/ro-apt-daily-upgrade.service
[Service]
ExecStartPre=-mount -o remount,rw /
ExecStartPost=-sh -c 'if [ "${NO_APT_REMOUNT:-no}" != yes ]; then mount -o remount,ro /; fi'